# nginx/permissions.conf
[program:permissions]
command=chown /var/run/docker.sock
autostart=true
autorestart=false
stderr_logfile=/var/log/permissions.err.log
stdout_logfile=/var/log/permissions.out.log
